# raspberrypi/etc/cron.d/pepoi2crrdtool
# 2018.5.4 add by pepo
MAILTO=""
*/5 * * * * root /usr/local/bin/pepoi2crrdtool
